# scripts/omens-daily-report.service
[Unit]
Description=omens daily FII report

[Service]
Type=oneshot
ExecStart=/home/kov/Projects/omens/scripts/daily-report.sh
