# scripts/omens-weekly-report.service
[Unit]
Description=omens weekly FII report

[Service]
Type=oneshot
ExecStart=/home/kov/Projects/omens/scripts/weekly-report.sh
